# qeos/tools/scripts/systemd/energyd.service
[Unit]
Description=QuantumEnergyOS Energy Daemon
Documentation=man:energyd(8)
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=qeos
Group=qeos
ExecStart=/usr/bin/energyd
Restart=on-failure
RestartSec=5s
LimitNOFILE=65536
Environment=RUST_LOG=info
Environment=RUST_BACKTRACE=1

[Install]
WantedBy=multi-user.target
